# docker_root/etc/services.d/uwsgi/run
#!/usr/bin/execlineb -P
/usr/local/bin/uwsgi --emperor /etc/uwsgi/sites --logto /var/log/uwsgi.log
